# get_location.service
[Unit]
Description=Get location

[Service]
Type=simple
ExecStart=%h/bin/get_location

[Install]
WantedBy=default.target
